# templates/sysstat-collect.timer.j2
[Timer]
OnCalendar=
OnCalendar={{ sysstat__collect_timer_on_calendar }}
AccuracySec={{ __sysstat__collect_timer_accuracy_sec }}
