# templates/sysstat-collect.timer.j2
[Timer]
OnCalendar=
OnCalendar={{ sysstat_collect_timer_on_calendar }}
AccuracySec={{ __sysstat_collect_timer_accuracy_sec }}
